# apple-tv-automation.service
[Unit]
Description=Apple TV Automation web app
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/apple-tv-automation
Environment=APPLE_TV_APP_HOST=0.0.0.0
Environment=APPLE_TV_APP_PORT=8000
ExecStart=/home/pi/apple-tv-automation/.venv/bin/python /home/pi/apple-tv-automation/server.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
